# resources/supervisord.conf
[unix_http_server]
file = /var/run/supervisor.sock
chmod = 0777
username = user
password = 123

[inet_http_server]
port = *:9001
username = user
password = 123

[supervisord]
logfile = /var/log/supervisor/supervisord.log
pidfile = /var/run/supervisord.pid
childlogdir = /var/log/supervisor

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock

[include]
files = /etc/supervisor/conf.d/*.conf

[program:cat]
command = /usr/bin/cat