# resources/supervisord.conf
[unix_http_server]
file = /var/run/supervisor.sock
chmod = 0777

[inet_http_server]
port = *:9001
username = user
password = 123

[supervisord]
logfile = /tmp/supervisord.log
pidfile = /tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock

[include]
files = /etc/supervisor/conf.d/*.conf

[program:cat]
command = /usr/bin/cat
